# etc/cron.d/bdii-proxy
37 5,11,17,23 * * * root /usr/sbin/bdii-proxy /opt/bdii/etc/bdii.conf
